# systemd_resolved_exporter/uos-systemd-resolved-exporter.service
[Unit]
Description=systemd resolved request metrics
Requires=network-online.target
After=network-online.target

[Service]
Restart=on-failure
ExecStart=/usr/bin/systemd-resolved-exporter
ExecReload=/bin/kill -HUP $MAINPID
TimeoutStopSec=20s
SendSIGKILL=no

[Install]
WantedBy=multi-user.target
